# .config/systemd/user/change-bg.service
[Unit]
Description=Execute bgaction script

[Service]
Type=oneshot
ExecStart=HOME=%h %h/.config/scripts/change-bg.sh

